# vitos-v1/packages/vitos-monitor/debian/postinst
#!/bin/sh
set -e
case "$1" in
  configure)
    getent group vitos-mon >/dev/null || groupadd --system vitos-mon
    getent passwd vitos-mon >/dev/null || \
      useradd --system --gid vitos-mon --home /var/lib/vitos --shell /usr/sbin/nologin vitos-mon
    install -d -o vitos-mon -g vitos-mon -m 0750 /var/log/vitos /run/vitos /var/lib/vitos

    python3 -m venv /opt/vitos-ai
    /opt/vitos-ai/bin/pip install --quiet -e /usr/lib/vitos/ai
    ln -sf /opt/vitos-ai/bin/vitos-ai /usr/bin/vitos-ai
    /opt/vitos-ai/bin/pip install --quiet -e /usr/lib/vitos/cli
    ln -sf /opt/vitos-ai/bin/vitosctl /usr/bin/vitosctl

    if command -v ollama >/dev/null 2>&1 && [ -f /etc/vitos/Modelfile ]; then
      ollama create vitos-intent -f /etc/vitos/Modelfile || true
    fi

    systemctl daemon-reload
    systemctl enable vitos-busd vitos-bpf-exec vitos-bpf-net \
      vitos-shell-tap vitos-udev-tap vitos-fanotify-tap \
      ollama vitos-ai || true
    ;;
esac
#DEBHELPER#
exit 0
