# vitos-v1/live-build/config/includes.chroot/etc/xdg/autostart/install-vitos-notice.desktop
[Desktop Entry]
Type=Application
Name=VITOS Install Notice
Exec=sh -c 'sleep 8 && notify-send -i calamares "Welcome to VITOS Live" "To install VITOS to your hard drive, click Install VITOS in the applications menu." -t 15000 2>/dev/null || true'
Terminal=false
Hidden=false
OnlyShowIn=XFCE;
X-XFCE-Autostart-Override=true
