# vitos-v1/packages/vitos-hardening/etc/cron.d/vitos-hardening
# /etc/cron.d/vitos-hardening — quarterly debsecan + lynis audit
# 03:00 on the 1st day of every 3rd month (Jan, Apr, Jul, Oct)
0 3 1 1,4,7,10 * root /usr/lib/vitos/hardening/run-audit.sh > /var/log/vitos/hardening/last-run.log 2>&1
