# vitos-v1/packages/vitos-monitor/systemd/vitos-bpf-exec.service
[Unit]
Description=VITOS eBPF execve collector
After=vitos-busd.service
Requires=vitos-busd.service

[Service]
Type=simple
ExecStart=/usr/sbin/vitos-bpf-exec --mode=exec --bus=/run/vitos/bus.sock
Restart=on-failure
AmbientCapabilities=CAP_BPF CAP_PERFMON CAP_SYS_RESOURCE
NoNewPrivileges=yes

[Install]
WantedBy=multi-user.target
